# dist-scripts/meetlogindex.service
[Unit]
Description=Meetbot Log Index Fetch Service

[Service]
Type=oneshot
ExecStart=/usr/bin/meetlogindex
